# scripts/pfc_supervisor.conf
[program:pfc]
autostart=true
startretries=10
startsecs=3
command=/usr/bin/pfc
redirect_stderr=true
stdout_logfile=/var/log/pfc.log
stopsignal=INT
