# scripts/pfs_supervisor.conf
[program:pfs]
autostart=true
startretries=10
startsecs=3
command=/usr/bin/pfs  -c /etc/pureflash/s5afs.conf
redirect_stderr=true
stdout_logfile=/var/log/pfs.log
stopsignal=INT
